# deploy/launchd/com.moltbook.worker.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.moltbook.worker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd <PROJECT_DIR> && uv run python -m src.workers.scheduler</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string><LOG_DIR>/moltbook-worker.out.log</string>
  <key>StandardErrorPath</key>
  <string><LOG_DIR>/moltbook-worker.err.log</string>
</dict>
</plist>
